# deploy/systemd/bizgenciz-web.service
[Unit]
Description=BizGenciz Web (Next.js Standalone)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bizgenciz-social
Environment=NODE_ENV=production
EnvironmentFile=/opt/bizgenciz-social/.env.local
ExecStart=/usr/bin/env bash -lc 'PORT=${PORT:-3000} /usr/bin/node /opt/bizgenciz-social/scripts/node-watchdog-wrapper.js /usr/bin/node .next/standalone/server.js'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStartSec=30
StartLimitIntervalSec=60
StartLimitBurst=5
# Optional watchdog: enable and implement sd_notify in the app to use. Example:
# WatchdogSec=60
# Type=notify

[Install]
WantedBy=multi-user.target


